# ledmatrix-wifi-monitor.service
[Unit]
Description=LED Matrix WiFi Monitor Daemon
After=network.target
Wants=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/home/ledpi/LEDMatrix
ExecStart=/usr/bin/python3 /home/ledpi/LEDMatrix/wifi_monitor_daemon.py --interval 30
Restart=on-failure
RestartSec=10
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=ledmatrix-wifi-monitor

[Install]
WantedBy=multi-user.target

